# scripts/rootfs/services/mintkit.service
[Unit]
Description=MintKit Launcher
After=multi-user.target

[Service]
Type=simple
User=mintkit
WorkingDirectory=/home/mintkit
ExecStart=/usr/bin/python3 /home/mintkit/launcher/mintos.py
Restart=on-failure
RestartSec=3
# KMS/DRM display — requires vc4-fkms-v3d overlay + raspberrypi-kernel (provides /dev/dri/card0)
Environment=SDL_VIDEODRIVER=kmsdrm
Environment=PYTHONPATH=/home/mintkit
StandardInput=null
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target